# packaging/systemd/user/aurynk-udev-proxy.socket
[Unit]
Description=Socket for Aurynk udev proxy

[Socket]
ListenStream=%t/aurynk-udev.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
